# docker/crontab
# Run auto-scaling check every minute
* * * * * /app/scripts/auto-scale.sh
